# discord-webhooks.service
[Unit]
Description=Discord Webhook Scripts for AssettoServer
After=network.target

[Service]
Type=forking
User=acserver
WorkingDirectory=/home/acserver/server
ExecStart=/home/acserver/server/start_webhooks.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
